# etc/systemd/system/llmnr.service
[Unit]
Description=LLMNR responder

[Service]
ExecStart=/usr/local/bin/llmnr

[Install]
WantedBy=multi-user.target